# order-coffee.service
[Unit]
Description=Order Coffee - Sleep Prevention Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/order-coffee --port 20553
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=true
ProtectHome=true
ReadWritePaths=/tmp

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
